# tasks.cron
# do daily/weekly/monthly maintenance
# min   hour    day     month   weekday command
*/15      *       *       *       *       speedtest -f jsonl | grep result | /root/record.sh
